# crontab
*/5 * * * * checkNameserver
0 0 * * * cleanDomains